# failing-unit-echo/echo@.service

[Unit]
Description=Display failing unit name

[Service]
Type=oneshot
ExecStart=/bin/echo "parent service : %i"
