# failing-unit-echo/failing-unit.service
[Unit]
Description=Failing systemd unit
After=network-online.target
OnFailure=echo@%i.service


[Service]
Type=oneshot
ExecStart=/bin/false
